# com.riadanas.whisperflow.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.riadanas.whisperflow</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/riadanas/Desktop/Fable 5/VoiceBud-Local-Riad/.venv/bin/python</string>
        <string>-u</string>
        <string>main.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/riadanas/Desktop/Fable 5/VoiceBud-Local-Riad</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/riadanas/Library/Logs/whisperflow.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/riadanas/Library/Logs/whisperflow.log</string>
</dict>
</plist>
